# container-root/etc/services.d/plex/finish
#!/usr/bin/with-contenv bash

rm -f "/config/Plex Media Server/plexmediaserver.pid"
